# scripts/systemd/greyhourrp-health-check.service
[Unit]
Description=GreyHourRP full health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/greyhourrp-automation
WorkingDirectory=/opt/greyhourrp
ExecStart=/bin/bash -c '/opt/greyhourrp/scripts/systemd/run-with-alert.sh health-check /opt/greyhourrp/scripts/health-check-all.sh'
