# scripts/systemd/greyhourrp-health-check.timer
[Unit]
Description=Run GreyHourRP health check every 15 minutes

[Timer]
OnCalendar=*-*-* *:00/15:00
Persistent=true
RandomizedDelaySec=30
Unit=greyhourrp-health-check.service

[Install]
WantedBy=timers.target
