# scripts/systemd/greyhourrp-self-heal.timer
[Unit]
Description=Run GreyHourRP self-heal every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=true
RandomizedDelaySec=30
Unit=greyhourrp-self-heal.service

[Install]
WantedBy=timers.target
