# scripts/systemd/greyhourrp-content-backup.timer
[Unit]
Description=Run GreyHourRP content backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=90
Unit=greyhourrp-content-backup.service

[Install]
WantedBy=timers.target
